# launchagents/com.theo.kill-idle-agents.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.theo.kill-idle-agents</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/theo/Projects/mac-install/bin/kill-idle-agents.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>1800</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/kill-idle-agents.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/kill-idle-agents.log</string>
</dict>
</plist>
